# build/service/irssi/run
#!/bin/bash
export HOME=/home/irssi/
set -e
echo "Starting session" >> /var/log/irssi-sessions.log
chpst -u irssi tmux new-session -s irssi -d irssi 2>&1 >> /var/log/irssi-sessions.log
echo "Session started" >> /var/log/irssi-sessions.log

set +e
while [ $? -eq 0 ]
do
  sleep 1
  chpst -u irssi tmux has-session -t irssi
done

echo "Session ended" >> /var/log/irssi-sessions.log